# roles/flexo/files/sched-paccache.timer
[Unit]
Description=paccache clean scheduling.

[Timer]
OnCalendar=weekly
Persistent=true

[Install]
WantedBy=timers.target